# systemd/mihomo-gateway.service.template
[Unit]
Description=Mihomo LAN Proxy Gateway
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/scripts/up.sh
ExecStop=__PROJECT_DIR__/scripts/down.sh
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
